# root/etc/services.d/sickbeard/run
#!/usr/bin/with-contenv bash

exec s6-setuidgid mab python /app/sickbeard/SickBeard.py --datadir /config

